# launchd/com.apoorav.profileagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apoorav.profileagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/apooravrao/Desktop/Profile Agent/orchestrator/main.py</string>
    </array>

    <!-- Run once daily at 10:00am -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Fire the missed job when Mac wakes from sleep -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/apooravrao/Desktop/Profile Agent/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/apooravrao/Desktop/Profile Agent/logs/launchd.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/apooravrao/Desktop/Profile Agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/apooravrao</string>
    </dict>
</dict>
</plist>
